# systemd/obsidian-headless-sync.service
[Unit]
Description=Obsidian Headless Sync
After=network-online.target
Wants=network-online.target

[Service]
User=obsidian
ExecStart=/usr/bin/ob sync --path "/home/obsidian/Vaults/a self-hosted homelab (Sync)" --continuous
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
