# systemd/obsidian-mcp.service
[Unit]
Description=Obsidian MCP Server
Documentation=https://github.com/${MCP_REPO}
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=bobsidian
Group=bobsidian
WorkingDirectory=/home/bobsidian/obsidian-mcp
EnvironmentFile=/home/bobsidian/.config/obsidian-mcp/env
ExecStart=/home/bobsidian/.local/bin/uv run vault-mcp
Restart=on-failure
RestartSec=10s
StartLimitBurst=5
StartLimitIntervalSec=60s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=obsidian-mcp

[Install]
WantedBy=multi-user.target
