# systemd/obsidian-git-sync.service
[Unit]
Description=Vault Git Sync
After=obsidian-headless-sync.service

[Service]
Type=oneshot
User=obsidian
EnvironmentFile=/home/obsidian/.config/obsidian-sync/env
ExecStart=/home/obsidian/obsidian-sync/scripts/obsidian-git-sync.sh
